# ntpclient.service.in
[Unit]
Description=Simple NTP client
Documentation=man:ntpclient
Documentation=file:@DOCDIR@/README.md
Documentation=https://github.com/troglobit/ntpclient
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
ExecStart=@SBINDIR@/ntpclient -n pool.ntp.org

[Install]
WantedBy=multi-user.target
